# ubuntu/overlay/lib/udev/rules.d/61-usbdevice.rules
SUBSYSTEM=="udc",ACTION=="change",DRIVER=="configfs-gadget",RUN+="/usr/bin/usb-hotplug %E{DEVPATH}"
