# debian/overlay/lib/systemd/system/async.service
[Unit]
Description=enable ASYNC for Debian Display
Before=lightdm.service

[Service]
Type=forking
ExecStart=/etc/init.d/S10atomic_commit.sh start
ExecStop=/etc/init.d/S10atomic_commit.sh stop
ExecReload=/etc/init.d/S10atomic_commit.sh reload

[Install]
WantedBy=multi-user.target

