# openwrt/package/firefly/ff_net_forward/files/etc/init.d/ff_net_forward
#!/bin/sh /etc/rc.common

START=99

start()
{
    /usr/bin/ff_net_forward
}
