# debian/overlay/etc/profile.d/x11.sh
export DISPLAY=${DISPLAY:-:0}
export TERM=xterm
