# openwrt/package/firefly/ff_net_forward/files/etc/hotplug.d/iface/99-ff_net_forward
[ ifup = "$ACTION" ] && [ "$INTERFACE" == "wan" ] && /usr/bin/ff_net_forward 
