# buildroot/board/rockchip/rk312x/fs-overlay/etc/profile.d/weston_override.sh
export QT_GSTREAMER_WIDGET_VIDEOSINK=kmssink
export QT_GSTREAMER_WINDOW_VIDEOSINK=kmssink
